# services/check_mk/systemd/files/check_mk@.service
[Unit]
Description=check_mk per connection handling

[Service]
ExecStart=-/usr/bin/check_mk_agent
StandardInput=socket
