# crowdrise/rclocal/files/rc.local
salt-call state.highstate
salt-call state.highstate
salt-call state.highstate
